# systemd/phoenix-ide.socket
[Unit]
Description=Phoenix IDE Socket
Documentation=https://github.com/phoenix-ide/phoenix-ide

[Socket]
# Production port - socket stays open during service restarts
ListenStream=8031
# Disable Nagle's algorithm for lower latency (SSE, interactive)
NoDelay=true
# Allow connections to queue during restart
Backlog=128
# Keep socket open when service stops (for zero-downtime restarts)
# Note: Service will be restarted by systemd when Accept=false (default)

[Install]
WantedBy=sockets.target
