# systemd/phoenix-ide.service
[Unit]
Description=Phoenix IDE
Documentation=https://github.com/phoenix-ide/phoenix-ide
# Socket must be ready before service starts
Requires=phoenix-ide.socket
After=network.target phoenix-ide.socket

[Service]
Type=simple
User=exedev

# Environment variables (PHOENIX_PORT not needed - socket provides the listener)
Environment=PHOENIX_DB_PATH=/home/exedev/.phoenix-ide/prod.db
# Note: PHOENIX_VERSION and LLM_GATEWAY are set dynamically by dev.py during deploy

ExecStart=/opt/phoenix-ide/phoenix-ide
# SIGHUP triggers graceful shutdown; systemd restarts with same socket
ExecReload=/bin/kill -HUP $MAINPID

# Restart on failure or when reloaded (SIGHUP causes exit 0)
Restart=always
RestartSec=1

# Give connections time to drain during graceful shutdown
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
